# indiestats.crontab
# m h  dom mon dow   command
20 9 * * * /home/indiestats/indiestats.sh > /home/indiestats/indiestats.cron.log 2>&1
